# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:openclaw]
command=/bin/sh /launcher.sh gateway --allow-unconfigured
user=node
environment=HOME="/home/node",NPM_CONFIG_PREFIX="/home/node/.npm-global",PATH="/home/node/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
